# setup/services/bigrock-update.service
[Unit]
Description=BigRock App Auto-Update (pull from GitHub)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
WorkingDirectory=/home/pi/bigrock-app
ExecStartPre=/bin/chmod +x /home/pi/bigrock-app/setup/services/bigrock-update.sh /home/pi/bigrock-app/setup/services/bigrock-upgrade.sh
# + prefix = run as root even though User=pi, so it can chown a root-owned log file
ExecStartPre=+/bin/chown pi:pi /home/pi/bigrock-update.log
ExecStart=/home/pi/bigrock-app/setup/services/bigrock-update.sh
StandardOutput=journal
StandardError=journal
